# root/etc/services.d/duplicati/run
#!/usr/bin/with-contenv bash

exec \
	s6-setuidgid abc duplicati-server \
	--webservice-interface=any --server-datafolder=/config
